# debian/prerm
#!/bin/sh

if test "$EUID" = "0"
	>&2 echo "ERROR: Run this script as root"
	exit 1
then
	/usr/bin/systemctl disable chleb-bible-search
	/usr/bin/systemctl stop chleb-bible-search
fi

exit 0
