# debian/prerm
#!/bin/sh

set -e

invoke-rc.d chleb-bible-search stop || :

#DEBHELPER#

exit 0
