# debian/prerm
#!/bin/sh

/usr/bin/systemctl disable chleb-bible-search
/usr/bin/systemctl stop chleb-bible-search

exit 0
