# debian/prerm
#!/bin/sh

set -e

systemctl disable chleb-bible-search
systemctl stop chleb-bible-search

#DEBHELPER#

exit 0
