# debian/prerm
#!/bin/sh

set -e

invoke-rc.d chleb-bible-search stop || :

echo 'clearing bible translation cache'
rm -f /var/cache/chleb-bible-search/*.bin

#DEBHELPER#

exit 0
