# debian/etc/cron.monthly/chleb-bible-search
#!/bin/sh

/usr/share/chleb-bible-search/session-clean.sh -t 30
